# ports/finit/rc.local
#!/bin/sh
for i in 1 2 3 4;do
	setleds -D +num < /dev/tty$i
done
cpupower frequency-set -g performance
