# ports/smartmontools/.footprint
bin/smartctl|755|0|0|regular file
bin/smartd|755|0|0|regular file
bin/update-smart-drivedb|755|0|0|regular file
etc/smartd.conf|644|0|0|regular file
etc/init.d/smartd|755|0|0|regular file
share/doc/smartmontools/README|644|0|0|regular file
share/doc/smartmontools/COPYING|644|0|0|regular file
share/doc/smartmontools/examplescripts/Example4|755|0|0|regular file
share/doc/smartmontools/examplescripts/README|644|0|0|regular file
share/doc/smartmontools/examplescripts/Example3|755|0|0|regular file
share/doc/smartmontools/examplescripts/Example2|755|0|0|regular file
share/doc/smartmontools/examplescripts/Example5|755|0|0|regular file
share/doc/smartmontools/examplescripts/Example1|755|0|0|regular file
share/doc/smartmontools/examplescripts/Example6|755|0|0|regular file
share/doc/smartmontools/smartd.conf|644|0|0|regular file
share/doc/smartmontools/NEWS|644|0|0|regular file
share/doc/smartmontools/TODO|644|0|0|regular file
share/doc/smartmontools/AUTHORS|644|0|0|regular file
share/doc/smartmontools/ChangeLog|644|0|0|regular file
share/doc/smartmontools/ChangeLog-5.0-6.0|644|0|0|regular file
share/doc/smartmontools/INSTALL|644|0|0|regular file
share/smartmontools/drivedb.h|644|0|0|regular file
share/smartmontools/smartd_warning.sh|755|0|0|regular file
share/man/man8/update-smart-drivedb.8.gz|644|0|0|regular file
share/man/man8/smartctl.8.gz|644|0|0|regular file
share/man/man8/smartd.8.gz|644|0|0|regular file
share/man/man5/smartd.conf.5.gz|644|0|0|regular file
share/man/man5|755|0|0|directory
share/man/man8|755|0|0|directory
share/man|755|0|0|directory
share/smartmontools/smartd_warning.d|755|0|0|directory
share/smartmontools|755|0|0|directory
share/doc/smartmontools/examplescripts|755|0|0|directory
share/doc/smartmontools|755|0|0|directory
share/doc|755|0|0|directory
share|755|0|0|directory
etc/init.d|755|0|0|directory
etc|755|0|0|directory
bin|755|0|0|directory
